# root/etc/s6-overlay/s6-rc.d/vhclient/finish
#!/usr/bin/with-contenv bash

mkdir -p /config/log/virtualhere
chown -R abc:abc /config/log/virtualhere
echo "[$(date)] VirtualHere client stopped" >> /config/log/virtualhere/vhclient.log